# tabulate/rootfs/etc/s6-overlay/s6-rc.d/nginx/run
#!/command/with-contenv sh
# ==============================================================================
# Start nginx reverse proxy
# ==============================================================================

echo "[tabulate] Starting nginx..."

exec nginx -g "daemon off;" -c /etc/nginx/nginx.conf
