# tabulate/rootfs/etc/s6-overlay/s6-rc.d/tabulate/finish
#!/command/with-contenv sh
# ==============================================================================
# Log exit of Tabulate backend
# ==============================================================================

echo "[tabulate] Tabulate backend stopped (exit code: ${1})"

if [ "${1}" -ne 0 ] && [ "${1}" -ne 256 ]; then
    echo "[tabulate] Backend crashed, restarting..."
fi
