# tabulate/rootfs/etc/s6-overlay/s6-rc.d/tabulate/run
#!/command/with-contenv sh
# ==============================================================================
# Start Tabulate backend (FastAPI + uvicorn)
# ==============================================================================

# Read config from HA add-on options
ANTHROPIC_API_KEY=$(jq -r '.anthropic_api_key' /data/options.json)
VERBOSE=$(jq -r '.verbose_logging // false' /data/options.json)

export ANTHROPIC_API_KEY
export DB_PATH="/data/tabulate.db"
export IMAGE_DIR="/data/images"

# Set log level based on addon option
if [ "$VERBOSE" = "true" ]; then
    export LOG_LEVEL="DEBUG"
    echo "[tabulate] Starting Tabulate backend (verbose logging)..."
    UVICORN_ARGS=""
else
    export LOG_LEVEL="INFO"
    echo "[tabulate] Starting Tabulate backend..."
    UVICORN_ARGS="--no-access-log"
fi

cd /app || exit 1

# shellcheck disable=SC2086
exec uvicorn main:app \
    --host 127.0.0.1 \
    --port 8000 \
    $UVICORN_ARGS
